# yggdrasil/deploy/munin/tts-server/yggdrasil-tts.service
[Unit]
Description=Yggdrasil Kokoro TTS Server
After=network.target
Before=yggdrasil-odin.service

[Service]
Type=simple
User=yggdrasil
Group=yggdrasil
SupplementaryGroups=render
WorkingDirectory=/opt/yggdrasil/deploy/munin/tts-server
ExecStart=/opt/yggdrasil/deploy/munin/tts-server/.venv/bin/python3 server.py
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
Environment=PYTHONUNBUFFERED=1
Environment=TTS_DEVICE=NPU

[Install]
WantedBy=multi-user.target
