# yggdrasil/deploy/systemd/claude-config-sync.service.example
[Unit]
Description=Claude Code Config Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=your-user
# Uses the user's SSH key for auth to the remote server
Environment=HOME=/home/your-user
EnvironmentFile=-/opt/yggdrasil/.env
ExecStart=/path/to/yggdrasil/deploy/workstation/claude-config-sync.sh sync
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
