# yggdrasil/deploy/systemd/claude-config-sync.timer
[Unit]
Description=Claude Code Config Sync Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
