# yggdrasil/deploy/systemd/yggdrasil-gitea-sync.service
[Unit]
Description=Yggdrasil Gitea Repo Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=your-user
Group=your-user
Environment=GITEA_TOKEN=$GITEA_TOKEN
Environment=GITEA_URL=http://<gitea-ip>:3000
Environment=REPO_DIR=/home/$USER/repos
Environment=GITEA_SSH_PORT=22
Environment=SKIP_REPOS=Crucible
ExecStart=/opt/yggdrasil/bin/gitea-sync.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
