# yggdrasil/deploy/systemd/yggdrasil-e2e-daily.service.example
[Unit]
Description=Yggdrasil daily E2E smoke — runs pytest tests-e2e against the live fleet
After=network-online.target yggdrasil-odin.service yggdrasil-mimir.service
Wants=network-online.target

# Sprint 066: cron wrapper now drives `pytest tests-e2e -m "not destructive
# and not slow"` instead of the legacy e2e-live.sh shell smoke. Pre-flight:
#   /opt/yggdrasil/tests-e2e/.venv/bin/pytest must exist and the suite must
#   be deployed to /opt/yggdrasil/tests-e2e. See scripts/smoke/e2e-cron-wrapper.sh
#   for the venv setup instructions printed on first failure.

[Service]
Type=oneshot
User=yggdrasil
Group=yggdrasil
WorkingDirectory=/opt/yggdrasil
EnvironmentFile=/opt/yggdrasil/.env
ExecStart=/opt/yggdrasil/scripts/smoke/e2e-cron-wrapper.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=20min
