# yggdrasil/deploy/systemd/yggdrasil-mcp-remote.service.example
[Unit]
Description=Yggdrasil MCP Remote Server (Streamable HTTP)
After=network.target yggdrasil-odin.service yggdrasil-mimir.service
Wants=yggdrasil-odin.service yggdrasil-mimir.service

[Service]
Type=simple
User=yggdrasil
ExecStart=/opt/yggdrasil/bin/ygg-mcp-remote --config /etc/yggdrasil/mcp-remote/config.yaml --bind 0.0.0.0:9093
Restart=on-failure
RestartSec=5
Environment="RUST_LOG=info"
# Set HA_TOKEN via a systemd drop-in or EnvironmentFile:
#   /etc/yggdrasil/mcp-remote/mcp-remote.env  →  HA_TOKEN=your-long-lived-token
# EnvironmentFile=/etc/yggdrasil/mcp-remote/mcp-remote.env
Environment="HA_TOKEN="

[Install]
WantedBy=multi-user.target
