# yggdrasil/deploy/systemd/yggdrasil-mcp-remote.service
[Unit]
Description=Yggdrasil MCP Remote Server (Streamable HTTP)
After=network.target yggdrasil-odin.service yggdrasil-mimir.service
Wants=yggdrasil-odin.service yggdrasil-mimir.service

[Service]
Type=simple
User=jhernandez
ExecStart=/opt/yggdrasil/bin/ygg-mcp-remote --config /etc/yggdrasil/mcp-remote/config.yaml --bind 0.0.0.0:9093
Restart=on-failure
RestartSec=5
Environment="RUST_LOG=info"
Environment="HA_TOKEN="

[Install]
WantedBy=multi-user.target
