# yggdrasil/deploy/munin/minicpm-server/yggdrasil-minicpm.service
[Unit]
Description=Yggdrasil MiniCPM-o Voice Server (ROCm GPU)
After=network.target
Before=yggdrasil-odin.service
OnFailure=yggdrasil-stt.service

[Service]
Type=simple
User=yggdrasil
WorkingDirectory=/opt/yggdrasil/deploy/munin/minicpm-server
ExecStart=/bin/bash start.sh
Restart=on-failure
RestartSec=30
TimeoutStartSec=600
Environment=ROCR_VISIBLE_DEVICES=0
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
