# yggdrasil/deploy/hugin/kokoro-tts/yggdrasil-kokoro-tts.service
[Unit]
Description=Yggdrasil Kokoro TTS Server
After=network.target

[Service]
Type=simple
User=jhernandez
WorkingDirectory=/opt/yggdrasil/deploy/kokoro-tts
ExecStart=/opt/yggdrasil/deploy/kokoro-tts/.venv/bin/python3 server.py --port 9098
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
