# yggdrasil/deploy/munin/yggdrasil-ollama-ipex.service
[Unit]
Description=Yggdrasil IPEX-LLM Ollama (Intel Arc iGPU)
Documentation=https://github.com/intel/ipex-llm
After=docker.service
Requires=docker.service
Before=yggdrasil-odin.service yggdrasil-mimir.service

[Service]
Type=simple
WorkingDirectory=/opt/yggdrasil/deploy/munin
ExecStartPre=/usr/bin/docker compose -f docker-compose.ipex-ollama.yml pull
ExecStart=/usr/bin/docker compose -f docker-compose.ipex-ollama.yml up --remove-orphans
ExecStop=/usr/bin/docker compose -f docker-compose.ipex-ollama.yml down
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
